# src/bin/dbab
#! /bin/sh
# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@";
    # if there is a /lib/init/init-d-script, then source it
    [ -s /lib/init/init-d-script ] &&
	INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
    # else, source my private copy to ease backporting
    [ -e /lib/init/init-d-script ] ||
	INIT_D_SCRIPT_SOURCED=true . /lib/init/dbab-init-d-script
fi
### BEGIN INIT INFO
# Provides:          dbab
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: dbab
# Description:       dnsmasq-based ad-blocking using pixelserv.
### END INIT INFO

DESC="dnsmasq-based ad-blocker"
DAEMON=/usr/sbin/dbab-svr

START_ARGS="--background --make-pidfile"
